# apiok-stash.conf
[program:apiok-stash]
directory=/opt/apiok-stash/current
command=/opt/apiok-stash/current/apiok-stash -f /opt/apiok-stash/current/etc/config.yaml
autostart=true
autorestart=true
startsecs=3
startretries=10
stopwaitsecs=30
stopsignal=TERM
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/opt/apiok-stash/current/stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
priority=999
